# roles/tone9_runtime/templates/quadlet/tone9.target.j2
[Unit]
Description=tone9 generated outline jobs
Wants=tone9-gen.service
After=tone9-gen.service

[Install]
WantedBy=default.target
